# etc/supervisord.d/kafka-connect.conf
[program:kafka-connect]
command=/usr/local/bin/start-kafka-connect.sh
autostart=%(ENV_ENABLE_BROKER)s
autorestart=true
startsecs=10
stdout_logfile=/var/log/%(program_name)s.log
redirect_stderr=true
